# infra/portfolio-contact.service
[Unit]
Description=Portfolio contact form API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=portfolio-contact
Group=portfolio-contact
WorkingDirectory=/opt/portfolio-contact/current
EnvironmentFile=/etc/portfolio-contact.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-contact

[Install]
WantedBy=multi-user.target
